# presto/scripts/it.presto.queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>it.presto.queue</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/francescogigante/Library/Application Support/Herd/bin/php</string>
        <string>/Users/francescogigante/Presto.it/presto/artisan</string>
        <string>queue:work</string>
        <string>--tries=3</string>
        <string>--sleep=3</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/francescogigante/Presto.it/presto</string>

    <key>StandardOutPath</key>
    <string>/tmp/presto-queue.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/presto-queue-error.log</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>
</dict>
</plist>
